# meta-custom/recipes-connectivity/wifi-config/files/wifi-init.service
[Unit]
Description=WiFi Interface Initialization
Before=wpa_supplicant@wlan0.service
After=sys-subsystem-net-devices-wlan0.device
Wants=sys-subsystem-net-devices-wlan0.device

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -c "ip link set wlan0 up"
ExecStart=/bin/sh -c "modprobe -r brcmfmac && sleep 2 && modprobe brcmfmac"

[Install]
WantedBy=multi-user.target