# meta-custom/recipes-connectivity/wifi-config/files/wifi-connect.service
[Unit]
Description=WiFi Connection Helper
After=wpa_supplicant@wlan0.service dhcpcd@wlan0.service
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/wifi-connect.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target